# systemd/ai-stack.target
[Unit]
Description=Sovereign Edge AI Stack
Documentation=https://github.com/t-timms/sovereign-edge
After=network-online.target ollama.service
Wants=network-online.target ollama.service voice.service

[Install]
WantedBy=multi-user.target
