# systemd/telegram-bot.service
[Unit]
Description=Sovereign Edge Telegram Bot
Documentation=https://github.com/YOUR_GITHUB_USERNAME/sovereign-edge
After=network-online.target ollama.service
Wants=network-online.target
PartOf=ai-stack.target

# ── Deployment Variables ───────────────────────────────────────────────────────
# Before enabling this service, set these to match your environment:
#   DEPLOY_USER  — the Linux user that owns the repo (e.g. your login username)
#   DEPLOY_ROOT  — absolute path to the sovereign-edge repo root
#
# Quick setup (run as root):
#   sed -i 's|<DEPLOY_USER>|YOUR_USER|g; s|<DEPLOY_ROOT>|/home/YOUR_USER/sovereign-edge|g' \
#       /etc/systemd/system/telegram-bot.service
#   systemctl daemon-reload

[Service]
Type=simple
User=<DEPLOY_USER>
WorkingDirectory=<DEPLOY_ROOT>
Environment=SOPS_AGE_KEY_FILE=<DEPLOY_ROOT>/../.config/sops/age/keys.txt
EnvironmentFile=-<DEPLOY_ROOT>/secrets/.env.decrypted

# Decrypt secrets — set -euo pipefail so sops failure aborts before sed writes
ExecStartPre=/bin/bash -c 'set -euo pipefail; \
    sops --decrypt <DEPLOY_ROOT>/secrets/env.yaml \
    | sed "s/: /=/" \
    > <DEPLOY_ROOT>/secrets/.env.decrypted \
    && chmod 600 <DEPLOY_ROOT>/secrets/.env.decrypted'

ExecStart=<DEPLOY_ROOT>/.venv/bin/python -m telegram_bot.bot
ExecStop=/bin/kill -TERM $MAINPID

Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sovereign-edge-telegram

# Resource limits — tune to your hardware
MemoryMax=2G
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=<DEPLOY_ROOT>

[Install]
WantedBy=ai-stack.target
